# overlay/etc/local.d/kernel-modules.start
#!/bin/sh
/sbin/depmod -a
modprobe btrfs
modprobe overlay
modprobe nfs
modprobe rpcsec_gss_krb5
