# overlay/etc/local.d/machine-id.start
#!/bin/sh
scw-metadata --cached ID > /etc/machine-id
