# overlay/etc/init.d/cgroupfs-volumes
#!/sbin/openrc-run
# Copyright (c) 2017 Marc Bachmann <marc.brookman@gmail.com>
# Released under the MIT license.

description="Mounts cgroupfs volumes for docker"

depend() {
    need net localmount
    before docker
}

start() {
    ebegin "Mount cgroup volumes"
    /usr/local/bin/cgroupfs-mount
    eend $?
}

stop() {
    ebegin "Unmount cgroup volumes"
    /usr/local/bin/cgroupfs-umount
    eend $?
}
