# overlay/etc/init.d/fluent-bit-system
#!/sbin/openrc-run
# Copyright (c) 2017 Marc Bachmann <marc.brookman@gmail.com>
# Released under the MIT license.

description="fluent-bit system log forwarder"
pidfile="/run/${RC_SVCNAME}.pid"
command="/usr/bin/env"
command_args="HOST_IP=${HOST_IP} HOST_NAME=${HOST_NAME} /usr/bin/fluent-bit -q -c /etc/fluent-bit/system.conf"
start_stop_daemon_args="--stdout /var/log/messages --stderr /var/log/messages"
command_background="true"

depend() {
    need hostname localmount
    before docker
    provide logger
}
