# overlay/etc/init.d/prometheus-openrc-exporter
#!/sbin/openrc-run
# Copyright (c) 2017 Marc Bachmann <marc.brookman@gmail.com>
# Released under the MIT license.

description="Export openrc service status using prometheus grok exporter"
pidfile="/run/${RC_SVCNAME}.pid"
command=/usr/local/bin/start-prometheus-openrc-exporter
start_stop_daemon_args="--stdout /var/log/${RC_SVCNAME} --stderr /var/log/${RC_SVCNAME}"
command_background="true"

depend() {
    need net root
    after docker
}
