# overlay/etc/init.d/prometheus-process-exporter
#!/sbin/openrc-run
# Copyright 2017 Marc Bachmann <marc.brookman@gmail.com>
# Released under the MIT license.

description="Prometheus process exporter"
pidfile="/run/${RC_SVCNAME}.pid"
command=/usr/local/sbin/process_exporter
command_args="-config.path /etc/prometheus/process_exporter.yml -web.listen-address 0.0.0.0:9256"
start_stop_daemon_args="--stdout /var/log/${RC_SVCNAME}.log --stderr /var/log/${RC_SVCNAME}.log"
command_background="true"

depend() {
  need net root
  before docker
}
