# overlay/etc/init.d/prometheus-node-exporter
#!/sbin/openrc-run
# Copyright (c) 2017 Marc Bachmann <marc.brookman@gmail.com>
# Released under the MIT license.

description="Prometheus node exporter"
pidfile="/run/${RC_SVCNAME}.pid"
command=/bin/node_exporter
command_args="-web.listen-address 172.17.0.1:9100"
start_stop_daemon_args="--stdout /var/log/${RC_SVCNAME} --stderr /var/log/${RC_SVCNAME}"
command_background="true"

depend() {
    need net root
    after docker
}
